# rootfs/etc/s6-overlay/s6-rc.d/init-dnsmasq/run
#!/command/execlineb -P

foreground { dnsmasq --version }
dnsmasq --test
