# rootfs/etc/s6-overlay/s6-rc.d/svc-dnsmasq/run
#!/command/execlineb -P

with-contenv
s6-notifyoncheck -d -n 300 -w 1000 -c "nc -z localhost 53" \
    dnsmasq \
        --keep-in-foreground \
        --log-facility=-
